# Myenv/nuevosiglo-homecast7582b0-secure-sdk-3.2rc5/Harvest/bSpec/nuevosiglo-hnc2200co-debug/inittab
# id:3:initdefault:
::sysinit:/etc/rcS 
::ctrlaltdel:/sbin/reboot -f
::respawn:/etc/gui_start
::respawn:/sbin/getty -n -i -l /bin/sh 115200 ttyS0 vt100
::shutdown:/sbin/reboot -f
::restart:/sbin/reboot -f

